# tests/to_config/system/d/dhcrelay4.service.d/00gentoo.conf
[Service]
Environment : "DHCRELAY_SERVERS="

